# etc/supervisor/conf.d/postfix.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:postfix]
command=/usr/local/bin/postfix.sh

[program:rsyslog]
command=/usr/sbin/rsyslogd -n -c3

[program:opendkim]
command=/usr/sbin/opendkim -f
